# kiss-tnc-bridge.service
[Unit]
Description=KISS TNC BLE-to-TCP Bridge
After=bluetooth.target network.target
Requires=bluetooth.target

[Service]
Type=simple
ExecStart=/usr/bin/kiss-tnc-bridge
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
